# tag-xenial/profile.d/40-rvm.sh
# -*- mode: sh; -*-

# Add RVM to PATH for scripting. Make sure this is the last PATH variable change.
[[ -d "$HOME/.rvm/bin" ]] && PATH=$PATH:$HOME/.rvm/bin
